# deploy/supervisor.conf
; ─────────────────────────────────────────────────────────────────────────────
; TienDo — Supervisor config cho Queue Worker
; File này đặt tại: /etc/supervisor/conf.d/tiendo-worker.conf
; Sau khi copy: sudo supervisorctl reread && sudo supervisorctl update
; ─────────────────────────────────────────────────────────────────────────────

[program:tiendo-worker]
command=/usr/bin/php8.2 /var/www/tiendo/artisan queue:work database \
    --queue=pdf-processing \
    --sleep=3 \
    --tries=3 \
    --backoff=30,60,120 \
    --timeout=120 \
    --max-jobs=500 \
    --memory=256

; Chạy với user www-data (không phải root)
user=www-data
directory=/var/www/tiendo

; Auto-start và auto-restart
autostart=true
autorestart=true
startretries=3
startsecs=3

; Số tiến trình — 1 worker đủ cho MVP
numprocs=1

; Logs
stdout_logfile=/var/log/supervisor/tiendo-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/tiendo-worker-error.log
stderr_logfile_maxbytes=5MB
redirect_stderr=false

; Killasgroup để đảm bảo child processes cũng bị dừng
killasgroup=true
stopasgroup=true
stopsignal=TERM
stopwaitsecs=30
